# etc/skel/.config/autostart/Startup Sound.desktop
[Desktop Entry]
Type=Application
Name=Startup Sound
Comment=Well, you get the idea.
Exec=start
